# stock-api.service
[Unit]
Description=Stock API (yfinance proxy)
After=network-online.target
Wants=network-online.target
# Never give up: don't let a burst of restarts auto-disable the unit.
StartLimitIntervalSec=0

[Service]
Type=simple
User=rozakos
WorkingDirectory=/home/rozakos/stock-api
EnvironmentFile=/home/rozakos/stock-api/.env
ExecStart=/home/rozakos/stock-api/.venv/bin/uvicorn main:app --host 127.0.0.1 --port 8001
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
